# deploy/ojas-backend.service
# Ojas backend — systemd unit.
#
# Place at  /etc/systemd/system/ojas-backend.service  then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable ojas-backend
#   sudo systemctl start ojas-backend
#   sudo systemctl status ojas-backend
#
# Logs:  journalctl -u ojas-backend -f

[Unit]
Description=Ojas — local agentic coding assistant (FastAPI)
After=network.target
Wants=network.target

[Service]
Type=simple
User=ojas
Group=ojas
WorkingDirectory=/opt/ojas

# Load .env BEFORE uvicorn starts so all the secrets are available. The file
# must be readable by the ojas user (chmod 600 + chown ojas:ojas).
EnvironmentFile=/opt/ojas/.env

# Use the venv's uvicorn so we don't depend on system Python paths.
#
# WebSocket tuning for production resilience:
#   --ws-ping-interval 20   send a ping every 20s so intermediate proxies
#                             (mobile carrier NATs, corporate firewalls)
#                             don't silently kill the idle socket
#   --ws-ping-timeout 20     close the connection if no pong within 20s
#                             (the client receives the close and its
#                             exponential-backoff reconnect kicks in)
#   --ws max_size 16MiB      cap individual message size; LLM tool
#                             results can occasionally be larger than the
#                             uvicorn default (2MiB)
ExecStart=/opt/ojas/.venv/bin/uvicorn server.app:app \
    --host 127.0.0.1 \
    --port 8765 \
    --log-level info \
    --proxy-headers \
    --forwarded-allow-ips 127.0.0.1 \
    --ws-ping-interval 20 \
    --ws-ping-timeout 20 \
    --ws-max-size 16777216

# Restart on crash with a small back-off so a tight crash loop doesn't
# spin up the box.
Restart=always
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=5

# Standard logging to journald — `journalctl -u ojas-backend` to read.
StandardOutput=journal
StandardError=journal

# Hardening — runtime sandbox so a compromised backend can't trivially
# escape and rm -rf /. Commented lines you can opt into for stricter
# setups; the defaults below are safe for a personal VM.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
# Ojas writes to /home/ojas/.agentic-rag and /home/ojas/ojas — those
# need to stay writable; allow them explicitly.
ReadWritePaths=/home/ojas

# File descriptor limit — plenty for many concurrent WebSocket clients.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
